# app_cron/crontab
# START CRON JOB
0 3 * * 0 python3 /app/cron.py
# END CRON JOB
